# demo-service/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:grpc]
command=/working/main -mode=grpc

[program:http]
command=/working/main -mode=http